# supervisor.conf
[supervisord]
nodaemon=true

[include]
files = /etc/supervisor/conf.d/*.conf

[program:webapp]
command=/bin/bash -c "apache2 -D FOREGROUND && /usr/sbin/sshd -D"
